# deploy/chat.conf
[program:chat]
directory= /usr/local/apps/chatter/
command = /usr/local/virtualenvs/chatter/env/bin/gunicorn -b 0.0.0.0:8080 chat.wsgi
user = root
