# deploy/chat.conf
[program:chat]
directory=/home/andrews/chatter/
command = /home/andrews/chatter/env/bin/gunicorn -b 0.0.0.0:8080 chat.wsgi
user = root
